# templates/trees/deb/flamerobin/etc/profile.d/flamerobin-path-env-exec.sh
#!/bin/bash



# Add Firebird and other custom paths to PATH
export PATH="/opt/flamerobin/bin:/opt/firebird/bin:/usr/local/bin:${PATH:-}"
